# scripts/com.portfolio.groww-reminder.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.portfolio.groww-reminder</string>
  <key>ProgramArguments</key>
  <array>
    <string>/path/to/portfolio/.venv/bin/python</string>
    <string>/path/to/portfolio/scripts/groww_subscription_reminder.py</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Month</key>
    <integer>6</integer>
    <key>Day</key>
    <integer>20</integer>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/path/to/portfolio</string>
</dict>
</plist>
